# tunnel-proxy/rootfs/etc/services.d/tunnel-proxy/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: 内网穿透代理
# 停止脚本
# ==============================================================================

bashio::log.info "正在停止内网穿透代理服务..."
